# root/etc/s6-overlay/s6-rc.d/init-beeper-config/run
#!/usr/bin/with-contenv bash

# Create Beeper config directories if they don't exist
mkdir -p /config/.config/BeeperTexts \
         /config/.config/beeper \
         /config/.local/share/beeper \
         /config/.pki/nssdb

# Fix ownership for abc user (PUID/PGID from environment, defaults to 1000)
chown -R abc:abc /config/.config \
                 /config/.local \
                 /config/.pki

echo "[init-beeper-config] Beeper config directories initialized with correct permissions"
