# root/.config/autostart/beeper.desktop
[Desktop Entry]
Type=Application
Name=Beeper
Exec=/usr/local/bin/beeper-launch
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
